# .crontab
SHELL=/bin/bash
PATH=/usr/local/bin:/usr/local/sbin:/sbin:/usr/sbin:/bin:/usr/bin:/usr/bin/X11
# each sunday at 00:05
#5 0 * * 0 cd /srv/data/stash/ && sudo fab backup > /dev/null
50 15 * * 3 cd /srv/data/stash/ && sudo fab backup > /dev/null
# -----------------------------------------------------------------------------
